# claude-telegram-bot.service
[Unit]
Description=Claude Code Telegram Bot
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/w/totoro-light/claude-telegram-bot
EnvironmentFile=/home/ubuntu/w/totoro-light/claude-telegram-bot/.env
ExecStart=/usr/bin/python3 /home/ubuntu/w/totoro-light/claude-telegram-bot/bot.py
Restart=always
RestartSec=5
KillMode=control-group
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-telegram-bot

[Install]
WantedBy=multi-user.target
